# rootfs/etc/init.d/rclone
#!/sbin/openrc-run

PIDFILE=/var/run/rclone.pid

depend() {
    before nfs
    provide portmap rpc.statd
}

start() {
    ebegin "Starting rclone"

    start-stop-daemon --start --background --pidfile "$PIDFILE" --make-pidfile --exec /usr/bin/rclone -- \
        mount remote: /mnt/rclone \
        --config=/etc/rclone/rclone.conf \
        --cache-dir=/var/rclone \
        --allow-other \
        --dir-cache-time=9999h \
        --log-file=/etc/rclone/rclone.log \
        --log-level=INFO \
        --umask=000 \
        --vfs-cache-mode=full \
        --vfs-fast-fingerprint \
        --vfs-cache-max-size=600G \
        --vfs-cache-max-age=9999h \
        --vfs-write-back=90m \
        --tpslimit 12 \
        --tpslimit-burst 0

    sleep 5

    exportfs -a

    eend $?
}

stop() {
    ebegin "Stopping rclone"

    start-stop-daemon --stop --quiet --pidfile $PIDFILE

    eend $?
}

status() {
    if mount | grep -q '/mnt/rclone'; then
        eend 0
    else
        eend 1
    fi
}
